# etc/profile.d/appmenu-gtk.sh
#!/bin/sh
export GTK_MODULES=unity-gtk-module:$GTK_MODULES
export UBUNTU_MENUPROXY=1
